# deploy/midiserver.service
[Unit]
Description=MIDI Visuals Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/midivisuals
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /var/www/midivisuals/server/index.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
